# .rMemory/config/com.hextrackr.memory-system.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hextrackr.memory-system</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Volumes/DATA/GitHub/HexTrackr/.rMemory/scripts/memory-system.sh</string>
        <string>start</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Volumes/DATA/GitHub/HexTrackr</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StartInterval</key>
    <integer>300</integer>
    
    <key>StandardOutPath</key>
    <string>/Volumes/DATA/GitHub/HexTrackr/.rMemory/logs/launchd-stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Volumes/DATA/GitHub/HexTrackr/.rMemory/logs/launchd-stderr.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/lbruton</string>
    </dict>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>AbandonProcessGroup</key>
    <true/>
</dict>
</plist>
