# docs/agents/rEngine-main/rCore/com.rengine.mcp-servers.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rengine.mcp-servers</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Volumes/DATA/GitHub/rEngine/integrated-mcp-manager.sh</string>
        <string>start</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Volumes/DATA/GitHub/rEngine/logs/mcp/launchd.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Volumes/DATA/GitHub/rEngine/logs/mcp/launchd.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Volumes/DATA/GitHub/rEngine</string>
</dict>
</plist>
