# app/scripts/com.hextrackr.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hextrackr.backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Volumes/DATA/GitHub/HexTrackr/app/scripts/hextrackr-backup.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Volumes/DATA/GitHub/HexTrackr/backups/hextrackr_backup.log</string>

    <key>StandardErrorPath</key>
    <string>/Volumes/DATA/GitHub/HexTrackr/backups/hextrackr_backup.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
